# supervisord.conf
[program:node]
command = node mapkiller_server.js
directory = /home/dotcloud/current